# spot-lyric-daemon/data/cn.spotlyric.Daemon.service.in
[Unit]
Description=Spot-Lyric daemon (Spotify lyrics overlay backend)
After=network-online.target

[Service]
Type=dbus
BusName=cn.spotlyric.Daemon
ExecStart=@bindir@/spot-lyric-daemon
Restart=on-failure
RestartSec=3

[Install]
WantedBy=default.target
